# edookit2summary.timer
# ABOUTME: Systemd timer for edookit summary gatherer.
# ABOUTME: Fires every 30 minutes; the script decides whether to work or just keepalive.

[Unit]
Description=Edookit school notification check

[Timer]
OnCalendar=*-*-* *:00,30:00
Persistent=true
Unit=edookit2summary.service

[Install]
WantedBy=timers.target
