# edookit2summary.timer
# ABOUTME: Systemd timer for edookit summary gatherer.
# ABOUTME: Fires daily at 15:00 to check for new school notifications.

[Unit]
Description=Edookit school notification check

[Timer]
OnCalendar=*-*-* 15:00:00
Persistent=true
Unit=edookit2summary.service

[Install]
WantedBy=timers.target
